# src/fty-certificate-generator.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=fty-certificate-generator service
After=network.target malamute.service fty-envvars.service fty-security-wallet.service
Requires=network.target malamute.service fty-envvars.service fty-security-wallet.service
Conflicts=shutdown.target
PartOf=bios.target

# config-path for certgen
ConditionPathExists=@datadir@/@PACKAGE@/

# We need the actual unix socket available to communicate with secw, not
# just the daemon running; otherwise currently we do not retry and abort :(
### src/certgen_certificate_generator_server.h:    static const std::string SECW_SOCKET_PATH = "/run/fty-security-wallet/secw.socket";
ConditionPathExists=@socketSecurityWallet@

[Service]
Type=simple
User=certificate-generator-daemon
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
EnvironmentFile=-/run/fty-envvars.env
ExecStart=@prefix@/bin/fty-certificate-generator
Restart=always

[Install]
WantedBy=bios.target
